# adm/bin/ptIInightly.cron
#!/bin/sh
# Version: @(#)ptIInightly.cron	1.5 12/01/98
if [ $# -ne 1 ]; then
    echo "$0: Usage: $0 email-address"
    exit 4	
fi
# Run ptIInightly and save the output

rm -f /users/ptII/adm/test/ptIInightly.log

# Run this with /bin/sh in case it is not executable
/bin/sh /users/ptII/adm/bin/ptIInightly $1 >/users/ptII/adm/test/ptIInightly.log 2>&1

cd /users/ptII/ptweb; cvs update -d -P > /users/ptII/logs/cvsupdate`date  '+%m%d'` 2>&1
